# debian/cron.d
MAILTO="peter.howe@emc.com"

*/1 * * * * root [ -x /usr/sbin/nullunitIterate ] && /usr/sbin/nullunitIterate
